# rendu/crontab
*/10 * * * * bash /home/rahmed/monitoring.sh | wall